# 10.0.0.8/supervisord.conf
[supervisord]
nodaemon=true

[program:mq_web]
command=/usr/local/bin/mq.sh 

[program:iib_manage]
user=iibuser
command=/usr/local/bin/iib_manage.sh
